# meta-rpi/recipes-headunit/headunit/files/headunit_start.service
[Unit]
Description=headunit
After=multi-user.target

[Service]
ExecStart=/usr/bin/headunit

[Install]
WantedBy=multi-user.target

